# components/dex/templates/dex-tunnel.service.j2
# dex SSH tunnel — user-level systemd unit on Linux clients.
#
# Forwards two ports on this machine through the GPU host's WSL SSH
# listener to their remote counterparts:
#   - {{ dex_local_port }} → {{ dex_remote_port }}: ollama (embed + chat)
#   - {{ dex_rerank_local_port }} → {{ dex_rerank_remote_port }}: Cohere-compatible rerank service (TEI/vLLM)
# ExitOnForwardFailure means either port unbound on the remote takes
# the tunnel down — both endpoints must be reachable on the server.
#
# Why an SSH tunnel rather than plain LAN HTTP:
#   - reuses existing SSH auth (no new bearer token to issue + rotate)
#   - keeps the embedding endpoint off the LAN (semantic embeddings of
#     your codebase are source-derived; treat them as sensitive)
#   - avoids a new Windows-FW + Hyper-V-FW rule pair on the GPU host
#
# Requires ~/.ssh/config to define `Host {{ dex_remote_host }}` with
# the right port + IdentityFile.
#
# ServerAlive* + ExitOnForwardFailure ensure a half-broken tunnel exits
# cleanly so systemd can restart it.

[Unit]
Description=dex SSH tunnel to {{ dex_remote_host }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ssh \
  -N \
  -o ExitOnForwardFailure=yes \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o StrictHostKeyChecking=accept-new \
  -L 127.0.0.1:{{ dex_local_port }}:127.0.0.1:{{ dex_remote_port }} \
  -L 127.0.0.1:{{ dex_rerank_local_port }}:127.0.0.1:{{ dex_rerank_remote_port }} \
  {{ dex_remote_host }}
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
