# components/dex/templates/dex-watch@.service.j2
# dex auto-reindex — systemd template unit (user scope).
#
# Instance name is the project root, systemd-escape'd. Enable per project:
#
#   systemctl --user enable --now \
#     dex-watch@$(systemd-escape /home/aleh/projects/dotfiles).service
#
# Disable:
#
#   systemctl --user disable --now \
#     dex-watch@$(systemd-escape /home/aleh/projects/dotfiles).service
#
# The watcher does a debounced incremental re-index on each save burst.
# It depends on the embedding endpoint being reachable (locally via the
# tunnel on clients, or 127.0.0.1 on the server). Restart=always covers
# the brief windows when the tunnel re-connects.

[Unit]
Description=dex watch %I
# Linux path escape: %I is the unescaped instance name (the project path).
# Wants/After dex-tunnel only on clients — harmless on the server
# where the unit doesn't exist (systemd ignores missing Wants).
Wants=dex-tunnel.service
After=dex-tunnel.service
# StartLimit* moved to [Unit] in systemd 230+ — [Service] is the legacy
# placement and is ignored with a warning. Tolerate a few restarts in
# 60 s, then back off.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# `%h` is systemd's user-home specifier (expanded at unit-load time).
# We can't use {{ dex_bin_path }} directly because it now resolves
# to the tilde-prefixed form (`~/.local/bin/dex`), and systemd does
# NOT expand `~` in ExecStart= — the unit would load with bad-setting.
# Embed-only: `dex watch` indexes but never generates, so the
# CHAT_* vars from the MCP registration in index.yml are intentionally
# omitted. EMBED_* still has to be set — ollama's model ref + non-default
# port can't be left to the binary's compiled-in defaults.
Environment=DEX_EMBED_URL=http://127.0.0.1:{{ dex_local_port }}
Environment=DEX_EMBED_MODEL={{ dex_embed_model }}
Environment=DEX_EMBED_BATCH={{ dex_embed_batch }}
ExecStart=%h/.local/bin/dex watch -v %I
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
