# components/dex/templates/dex-summary-tunnel.service.j2
# dex summary tunnel — user-level systemd unit on the GPU server.
#
# Forwards the remote summary worker's ollama port to localhost so
# DEX_SUMMARY_URL (127.0.0.1:{{ dex_summary_port }}) reaches mini_pc's
# ollama (127.0.0.1:11434) over SSH without exposing it on the LAN.
#
# Mirrors the design of dex-tunnel.service on clients: SSH auth (no bearer
# token), ExitOnForwardFailure for clean systemd restarts on half-broken
# tunnels, ServerAlive* for dead-connection detection.
#
# Managed by components/dex/index.yml — do not edit by hand.

[Unit]
Description=dex summary tunnel to {{ dex_summary_host }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ssh \
  -N \
  -o ExitOnForwardFailure=yes \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o StrictHostKeyChecking=accept-new \
  -L 127.0.0.1:{{ dex_summary_local_port }}:127.0.0.1:11434 \
  {{ dex_summary_host }}
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
